# houston.socket
[Socket]
ListenStream=10101

[Install]
WantedBy=sockets.target
