# houston.service
[Service]
ExecStart=/usr/bin/houston serve
